# debian/glance.prerm
#!/bin/sh                                                                                                                                  

set -e

case $1 in
        remove|purge)
                invoke-rc.d glance-api stop
                invoke-rc.d glance-registry stop
esac

#DEBHELPER#