# debian/glance-registry.prerm
#!/bin/sh

set -e

. /usr/share/debconf/confmodule

db_get glance-registry/configure_db
if [ "$RET" = "true" ]; then
    . /usr/share/dbconfig-common/dpkg/prerm
    dbc_go glance-registry $@
fi

case $1 in
    remove|purge)
	if [ -x /etc/init.d/glance-registry ]; then
	    invoke-rc.d glance-registry stop || true
	fi
	;;
esac

#DEBHELPER#
