# debian/glance-api.prerm
#!/bin/sh

set -e

case $1 in
	remove|purge)
		invoke-rc.d glance-api stop
esac

#DEBHELPER#
